# systemd/remer.us.path
[Unit]
Description=Watch for changes to remer.us site

[Path]
PathChanged=/usr/local/lib/node_modules/remer.us/package.json
Unit=remer.us-upgrade.service

[Install]
WantedBy=multi-user.target
